# ops/launchd/com.groovejones.paperclip.mini.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.groovejones.paperclip</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd '/Users/chipcarman/dev/paperclip' &amp;&amp; exec pnpm paperclipai run</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/chipcarman/dev/paperclip</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/chipcarman/Library/Logs/paperclip.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/chipcarman/Library/Logs/paperclip.err</string>
</dict>
</plist>
