# srv.service
[Unit]
Description=Paperclip live dev service
After=network.target

[Service]
Type=simple
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/paperclip
Environment=HOME=/home/exedev
Environment=PATH=/home/exedev/.local/bin:/home/exedev/node/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/bash -lc 'set -a; [ -f .env ] && . ./.env; set +a; exec /home/exedev/.local/bin/pnpm dev:once'
Restart=always
RestartSec=5
KillMode=control-group
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
