# contrib/macos-launchagent/io.paperclip.local.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.paperclip.local</string>
  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/paperclip-repo</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/ABSOLUTE/PATH/TO/paperclip-repo/cli/node_modules/tsx/dist/cli.mjs</string>
    <string>/ABSOLUTE/PATH/TO/paperclip-repo/cli/src/index.ts</string>
    <string>run</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/USERNAME</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PAPERCLIP_HOME</key>
    <string>/Users/USERNAME/.paperclip</string>
    <key>NODE_OPTIONS</key>
    <string>--max-old-space-size=8192</string>
    <key>PAPERCLIP_MANAGED_BY_LAUNCHD</key>
    <string>1</string>
    <key>PAPERCLIP_STRICT_LISTEN_PORT</key>
    <string>true</string>
    <key>PAPERCLIP_UI_DEV_MIDDLEWARE</key>
    <string>false</string>
    <key>SERVE_UI</key>
    <string>true</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/USERNAME/.paperclip/launchd.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/USERNAME/.paperclip/launchd.stderr.log</string>
</dict>
</plist>
